# scripts/assistext-ollama-loader.service
[Unit]
Description=Assistext Ollama Model Loader
After=ollama.service network.target
Requires=ollama.service

[Service]
Type=oneshot
ExecStart=/bin/bash /path/to/Assistext/scripts/preload_ollama_model.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
